# root/etc/s6-overlay/s6-rc.d/serpbear/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

echo "Starting serpbear..."

# Set default values
PORT=${PORT:-3000}
NODE_ENV=${NODE_ENV:-production}

# Export environment variables for Node.js
export NODE_ENV
export PORT

echo "Starting serpbear on port ${PORT}..."

# Change to application directory
cd /app || exit 1

# Run database migrations first (like original image)
echo "Running database migrations..."
npx sequelize-cli db:migrate --env production 2>/dev/null || true

# Start the application directly
echo "Starting Node.js server..."
exec node server.js