# root/etc/s6-overlay/s6-rc.d/init-serpbear-config/up
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

echo "Configuring serpbear..."

# Create basic directories
mkdir -p /config /data 2>/dev/null || true

echo "✓ serpbear configuration complete"