# root/etc/s6-overlay/s6-rc.d/init-custom-files/up
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

echo "Custom files init (LinuxServer.io standard)"
echo "✓ Custom files init complete"